# conf/supervisor.conf
# -*- conf -*-

[program:dashbeat]
command=node /srv/sites/dashbeat/dist/server.js
user=ebower
group=web
stopsignal=QUIT
autostart=true
autorestart=true
stderr_logfile=syslog
stdout_logfile=syslog
environment=NODE_PORT="5001",NODE_ENV="production"
